# alpine/openrc/mtg/mtg.initd
#!/sbin/openrc-run

command="/usr/local/bin/mtg"
command_args="simple-run $BIND_TO $SECRET --timeout=$TIMEOUT --prefer-ip=$PREFER_IP --doh-ip=$DOH_IP"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
start_stop_daemon_args="--nicelevel -5"

depend() {
    need net
}
